# unl0kr-ask-password.path
#  SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Forward Password Requests to unl0kr Directory Watch
DefaultDependencies=no
Conflicts=shutdown.target emergency.service
After=plymouth-start.service
Before=paths.target shutdown.target cryptsetup.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes

